# data_collection_transfer/scheduling/collect_metric_timer.service
[Unit]
Description=Collect System Metrics Timer

[Timer]
OnUnitActiveSec=60s
OnBootSec=60s
Unit=collect_metrics.service

[Install]
WantedBy=timers.target
